# dotfiles/.config/systemd/user/pueue-daemon.service
[Unit]
Description=Pueue Daemon

[Service]
ExecStart=%h/.cargo/bin/pueued
Restart=always
RestartSec=15

[Install]
WantedBy=default.target
